# docker/root/etc/s6-overlay/s6-rc.d/svc-web/run
#!/command/with-contenv sh
# shellcheck shell=bash

cd /app/apps/web;
exec node server.js;
