# docker/root/etc/s6-overlay/s6-rc.d/svc-workers/run
#!/command/with-contenv sh
# shellcheck shell=bash

cd /app/apps/workers;
exec pnpm run start:prod;
