# docker/root/etc/s6-overlay/s6-rc.d/init-db-migration/run
#!/command/with-contenv sh
# shellcheck shell=bash

echo "Running db migration script";
cd /db_migrations;
exec node index.js;
